# extra/omnivore-exporter.timer
[Unit]
Description=Omnivore Exporter Timer

[Timer]
OnUnitActiveSec=60min 

[Install]
WantedBy=timers.target

